# support/ubuntu/pyprobe.conf
# pyprobe - Mini Probe for PRTG

description     "Monitoring daemon for PRTG"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn

exec /usr/local/bin/pyprobe
